# code-du-travail-bot.service
# Configuration systemd pour le bot Code du Travail
# Ce fichier sera utilisé par setup.sh pour créer le service

[Unit]
Description=Bot Telegram Code du Travail - Mistral 7B
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Code-du-Travail
Environment=PATH=/home/ubuntu/Code-du-Travail/venv/bin
ExecStart=/home/ubuntu/Code-du-Travail/venv/bin/python /home/ubuntu/Code-du-Travail/telegram_bot.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=code-du-travail-bot

# Limite de ressources (optionnel)
LimitMEMLOCK=infinity
LimitNOFILE=65536

# Variables d'environnement
EnvironmentFile=-/home/ubuntu/Code-du-Travail/.env

[Install]
WantedBy=multi-user.target
